# packages/xerotier-xim-agent/prerm
#!/bin/sh
# SPDX-License-Identifier: MIT
set -e

if [ "$1" = "remove" ] && command -v systemctl >/dev/null 2>&1; then
    systemctl --no-reload disable --now xerotier-xim-agent.service >/dev/null 2>&1 || true
fi
